# src/dispatcher/supervisord.conf
[supervisord]
nodaemon=true

[program:dispatcher]
command=python3 /tmp/dispatcher.py
